# mzansiedge.service
[Unit]
Description=MzansiEdge Telegram Bot
After=network.target

[Service]
Type=simple
User=paulsportsza
WorkingDirectory=/home/paulsportsza/bot
EnvironmentFile=/home/paulsportsza/bot/.env
ExecStart=/usr/bin/python3 /home/paulsportsza/bot/bot.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
